# systemd/streamdeploy-agent.service
[Unit]
Description=StreamDeploy Agent
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/streamdeploy-agent /etc/streamdeploy/agent.json
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
